# debian/metricsd.postinst
chown nobody:nogroup /var/log/metricsd
update-rc.d metricsd defaults
